# target/linux/ramips/mt7621/base-files/etc/init.d/set-irq-affinity
#!/bin/sh /etc/rc.common
#MTK-Wi-Fi

chmod 0775 /sbin/mtk-wifi
chmod 0775 /sbin/mtk-wifi.sh

START=99

start() {
	/sbin/mtk-wifi.sh &
}
